# src/packages/fff/fff-sysupgrade/files/etc/uci-defaults/99-fff-sysupgrade
# SPDX-License-Identifier: GPL-3.0-only
#
# Copyright 2017 Tim Niemeyer

cat > /etc/sysupgrade.conf <<-__EOF__
/etc/shadow
/etc/dropbear/dropbear_rsa_host_key
/etc/dropbear/dropbear_ecdsa_host_key
/etc/dropbear/dropbear_ed25519_host_key
/etc/dropbear/authorized_keys
/etc/network.config
/etc/config/fff
/etc/config/gateway
/etc/hoodfile
__EOF__

exit 0
