# bsp/default/root_file_system/etc/rc.local
#!/bin/sh
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

/usr/sbin/configurenetwork

touch /tmp/started

exit 0
