# src/packages/fff/fff-hoods/files/etc/uci-defaults/24c-fff-wXsta
#!/bin/sh

. /lib/functions/fff/wireless

# Set up wXsta
for radio in $(wifiListRadio); do
	# wXsta: We can use $freq here, as the initial radios are not set up with auto
	freq="$(wifiGetFreq $radio)"

	uci batch <<-__EOF__
		set network.configsta${freq}=interface
		set network.configsta${freq}.proto='static'

		set wireless.w${freq}sta='wifi-iface'
		set wireless.w${freq}sta.device='${radio}'
		set wireless.w${freq}sta.network='configsta${freq}'
		set wireless.w${freq}sta.ifname='w${freq}sta'
		set wireless.w${freq}sta.mode='sta'
		set wireless.w${freq}sta.ssid='config.franken.freifunk.net'
		set wireless.w${freq}sta.disabled='1'
	__EOF__
done

uci commit network
uci commit wireless

# vim: set noexpandtab:tabstop=4
