# src/packages/fff/fff-ra/files/etc/uci-defaults/90-fff-ra
uci batch <<EOF
	set dhcp.client=dhcp
	set dhcp.client.interface='client'
	set dhcp.client.ra='server'
	set dhcp.client.ra_default='2'
	set dhcp.client.ra_management='0'
	set dhcp.client.ra_useleasetime='1'
EOF

exit 0
