# src/packages/fff/fff-ra/files/etc/uci-defaults/90-fff-ra
#!/bin/sh

uci batch <<EOF
	set dhcp.mesh=dhcp
	set dhcp.mesh.interface='mesh'
	set dhcp.mesh.ra='server'
	set dhcp.mesh.ra_default='2'
	set dhcp.mesh.ra_management='0'
EOF

exit 0
