# src/packages/fff/fff-network/files/etc/uci-defaults/20-migrate-c7-portconfig
BOARD="$(uci get board.model.name)"

# only execute if running on a TP-Link Archer C7 v2
if [ "$BOARD" != "tplink,archer-c7-v2" ]; then
	exit 0
fi

# only execute if there is something to migrate
if [ ! -s /etc/network.config ]; then
	exit 0
fi

# modify WAN_PORTS to new single-interface configuration
sed -i '/WAN_PORTS/ s/6/0t/g' /etc/network.config
