# src/packages/fff/fff-network/files/etc/uci-defaults/20-remove-ports-override
# Remove port configuration residues from earlier firmware versions
sed -i -e '/^WAN_PORTS/d' -e '/^BATMAN_PORTS/d' -e '/^CLIENT_PORTS/d' /etc/network.config
