# src/packages/fff/fff-network/files/etc/uci-defaults/50-packet-steering
#!/bin/sh
# SPDX-License-Identifier: GPL-3.0-or-later

# Enable RPS
uci batch <<-__EOF__
	set network.globals=globals
	set network.globals.packet_steering='1'
__EOF__

uci -q commit network
