# src/packages/fff/fff-network/files/etc/uci-defaults/94-fff-nodewatcher-network
uci batch <<EOF
  add nodewatcher network
  set nodewatcher.@network[-1].mesh_interface='br-mesh'
  set nodewatcher.@network[-1].iface_blacklist='lo ifb0'
  set nodewatcher.@network[-1].ip_whitelist='br-mesh'
EOF
uci commit nodewatcher

exit 0
