# src/packages/fff/fff-support/files/etc/uci-defaults/80-banner-info
# Copyright 2018 Adrian Schmutzler
# License GPLv3

. /etc/firmware_release

echo "" >> /etc/banner
echo "> $(cat /tmp/sysinfo/model 2>/dev/null) @ $FIRMWARE_VERSION" >> /etc/banner
echo "" >> /etc/banner

exit 0
