# src/packages/fff/fff-support/files/etc/uci-defaults/80-banner-info
# SPDX-License-Identifier: GPL-3.0-only
#
# Copyright 2018 Adrian Schmutzler

. /etc/firmware_release

echo "" >> /etc/banner
echo "> $(cat /tmp/sysinfo/model 2>/dev/null) @ $FIRMWARE_VERSION" >> /etc/banner
echo "" >> /etc/banner

exit 0
