# src/packages/fff/fff-config/files/etc/uci-defaults/98-configure-fff
#!/bin/sh
# Copyright 2017 Adrian Schmutzler
# License GPLv3

touch /etc/config/fff

if ! uci -q get fff.system > /dev/null ; then
	uci set fff.system=fff
	uci set "fff.system.hostname=FFF"
	uci commit fff
fi

FFF_hostname=$(uci -q get "fff.system.hostname")
test -n "${FFF_hostname}" && uci -q set "system.@system[0].hostname=${FFF_hostname}"
test -n "${FFF_hostname}" && echo "${FFF_hostname}" > /proc/sys/kernel/hostname

uci commit system

poe_passthrough=$(uci -q get "fff.poe_passthrough.active")
if [ "$poe_passthrough" = "1" ] && [ -s /usr/lib/fff-support/activate_poe_passthrough.sh ] ; then
	/usr/lib/fff-support/activate_poe_passthrough.sh
fi
