# src/packages/fff/fff-fastd/files/etc/uci-defaults/55_fff-fastd
#!/bin/sh

/etc/init.d/fastd disable

>/etc/config/fastd

uci batch <<EOF
  set fastd.fff='fastd'
  set fastd.fff.enabled='1'
  set fastd.fff.config_peer_dir="/etc/fastd/fff/peers"
  set fastd.fff.syslog_level='warn'
  set fastd.fff.method='null'
  set fastd.fff.mode='tap'
  set fastd.fff.interface="fffVPN"
  set fastd.fff.mtu='1426'
  set fastd.fff.on_up="/etc/fastd/fff/up.sh"
  set fastd.fff.secure_handshakes='0'
  set fastd.fff.secret="generate"
EOF

[ ! -d /etc/fastd/fff ] &&  mkdir -p /etc/fastd/fff
ln -s /tmp/fastd_fff_peers /etc/fastd/fff/peers
echo "#!/bin/sh" > /etc/fastd/fff/up.sh
echo "ip link set up dev fffVPN" >> /etc/fastd/fff/up.sh
echo "echo enable > /sys/devices/virtual/net/fffVPN/batman_adv/no_rebroadcast" >> /etc/fastd/fff/up.sh
echo "batctl if add fffVPN" >> /etc/fastd/fff/up.sh
chmod +x /etc/fastd/fff/up.sh
