# src/packages/fff/fff-nodewatcher/files/etc/hotplug.d/iface/51-configurator
#!/bin/sh
[ "$ACTION" = "ifup" -a "$INTERFACE" = "wan" ] && {
	/usr/sbin/configurator
}
