# src/packages/fff/fff-simple-tc/files/etc/uci-defaults/98-configure-tc
#!/bin/sh

traffic_enabled=$(uci -q get "fff.trafficcontrol.enabled")
traffic_in=$(uci -q get "fff.trafficcontrol.limit_in")
traffic_out=$(uci -q get "fff.trafficcontrol.limit_out")

test -n "${traffic_enabled}" && uci set "simple-tc.example.enabled=${traffic_enabled}"
uci set "simple-tc.example.ifname=wan"
uci set "simple-tc.example.limit_ingress=${traffic_in}"
uci set "simple-tc.example.limit_egress=${traffic_out}"

uci commit simple-tc
