# src/packages/fff/fff-tunneldigger/files/etc/uci-defaults/55_fff-tunneldigger
/etc/init.d/tunneldigger disable
>/etc/config/tunneldigger
