# src/packages/fff/fff-node/files/etc/uci-defaults/23-network-node
#!/bin/sh
# SPDX-License-Identifier: GPL-3.0-or-later

. /lib/functions/system.sh
. /lib/functions/fff/network

BOARD="$(uci get board.model.name)"
. /etc/network.$BOARD

[ -n "$ROUTERMAC" ] || ROUTERMAC=$(get_mac_label)

prefix="fdff::/64"

# Set $prefix::1 as IP
addr1=$(owipcalc "$prefix" add ::1)
# Set $prefix::MAC as IP
addr2=$(owipcalc "$prefix" add "::$(ipMacSuffix "$ROUTERMAC")")
# Set $prefix::EUI64 as IP
addr3=$(owipcalc "$prefix" add "::$(ipEUISuffix "$ROUTERMAC")")

uci batch <<-__EOF__
	set network.globals.ula_prefix=$prefix

	add_list network.client.ip6addr="${addr1}"
	add_list network.client.ip6addr="${addr2}"
	add_list network.client.ip6addr="${addr3}"
	set network.client.proto=static

	set network.ethmesh='interface'
	set network.ethmesh.proto='batadv_hardif'
	set network.ethmesh.master='bat0'
__EOF__


uci -q commit network
