# src/packages/fff/fff-node/files/etc/uci-defaults/23-network-node
#!/bin/sh
# SPDX-License-Identifier: GPL-3.0-or-later

. /lib/functions/system.sh
. /lib/functions/fff/network
. /lib/functions/fff/cpuport

BOARD="$(uci get board.model.name)"
. /etc/network.$BOARD

[ -n "$ROUTERMAC" ] || ROUTERMAC=$(get_mac_label)

prefix="fdff::/64"

# Set $prefix::1 as IP
addr1=$(owipcalc "$prefix" add ::1)
# Set $prefix::MAC as IP
addr2=$(owipcalc "$prefix" add "::$(ipMacSuffix "$ROUTERMAC")")
# Set $prefix::EUI64 as IP
addr3=$(owipcalc "$prefix" add "::$(ipEUISuffix "$ROUTERMAC")")

# Add client network configuration
uci batch <<-__EOF__
	set network.globals.ula_prefix=$prefix

	set network.client='interface'
	set network.client.type='bridge'
	set network.client.auto='1'
	set network.client.proto=static
	set network.client.fff_clientif='1'

	add_list network.client.ip6addr="${addr1}"
	add_list network.client.ip6addr="${addr2}"
	add_list network.client.ip6addr="${addr3}"

	set network.ethmesh='interface'
	set network.ethmesh.proto='batadv_hardif'
	set network.ethmesh.master='bat0'
__EOF__


# Add cpuport to switch configuration
CLIENT_PORTS="$(get_cpu_port) $CLIENT_PORTS"
BATMAN_PORTS="$(get_cpu_port) $BATMAN_PORTS"
[ -n "$WAN_PORTS" ] && WAN_PORTS="$(get_cpu_port) $WAN_PORTS"

SWITCHHW=$(swconfig list | awk '{ print $4 }')
if [ "$DSA" = "1" ]; then
	uci batch <<-__EOF__
		set network.${SWITCHDEV}_1=bridge-vlan
		set network.${SWITCHDEV}_1.device=$SWITCHDEV
		set network.${SWITCHDEV}_1.vlan=1
		set network.${SWITCHDEV}_1.ports="$CLIENT_PORTS"

		set network.${SWITCHDEV}_2=bridge-vlan
		set network.${SWITCHDEV}_2.device=$SWITCHDEV
		set network.${SWITCHDEV}_2.vlan=2
		set network.${SWITCHDEV}_2.ports="$WAN_PORTS"

		set network.${SWITCHDEV}_3=bridge-vlan
		set network.${SWITCHDEV}_3.device=$SWITCHDEV
		set network.${SWITCHDEV}_3.vlan=3
		set network.${SWITCHDEV}_3.ports="$BATMAN_PORTS"
	__EOF__
elif [ -n "$SWITCHHW" ] && [ -z "$ONE_PORT" ] && [ -z "$TWO_PORT" ]; then
	uci batch <<-__EOF__
		set network.${SWITCHDEV}_1=switch_vlan
		set network.${SWITCHDEV}_1.device=$SWITCHHW
		set network.${SWITCHDEV}_1.vlan=1
		set network.${SWITCHDEV}_1.ports="$CLIENT_PORTS"

		set network.${SWITCHDEV}_2=switch_vlan
		set network.${SWITCHDEV}_2.device=$SWITCHHW
		set network.${SWITCHDEV}_2.vlan=2
		set network.${SWITCHDEV}_2.ports="$WAN_PORTS"

		set network.${SWITCHDEV}_3=switch_vlan
		set network.${SWITCHDEV}_3.device=$SWITCHHW
		set network.${SWITCHDEV}_3.vlan=3
		set network.${SWITCHDEV}_3.ports="$BATMAN_PORTS"
	__EOF__
fi

if [ -n "$SWITCHHW" ] || [ "$DSA" = "1" ]; then
	uci set network.client.ifname="$SWITCHDEV.1 bat0"
	uci set network.wan.ifname=$WANDEV.2
	uci set network.ethmesh.ifname="$SWITCHDEV.3"
fi

if [ "$WANDEV" != "$SWITCHDEV" ]; then
	# wan interface is not part of the onboard switch.
	# one- and twoport devices are handled by configurenetwork.
	uci set network.wan.ifname=$WANDEV
fi

if [ -n "$ETHMESHMAC" ]; then
	uci set network.ethmesh.macaddr=$ETHMESHMAC
fi

if [ -n "$ROUTERMAC" ]; then
	uci set network.client.macaddr=$ROUTERMAC
fi


uci batch <<-__EOF__
__EOF__


uci -q commit network
