# src/packages/fff/fff-layer3-config/files/etc/init.d/fff-layer3-update
#!/bin/sh /etc/rc.common

START=99

boot() {
	# Reapply gateway settings, if configuration succeeds
	if yes | configure-layer3 -c; then
		configure-layer3 -a
	else
		configure-layer3 -r
	fi

	# first we disable the init.d
	/etc/init.d/fff-layer3-update disable
	# we must delete the symlink manually
	rm -f /etc/rc.d/S99fff-layer3-update
}
