# src/packages/fff/fff-layer3-config/files/etc/uci-defaults/21-l3config-migrate-3-to-4
## this script migrates configs from version 2 to 3.
##
## with version 2 the switchport names of some devices changed
## due to the switch to DSA


BOARD="$(uci get board.model.name)"

# only migrate appropriate config versions
[ "$(uci -q get gateway.meta.config_version)" = "3" ] || exit 0


translate_ports() {
	local vlan="$1"
	local oldports="$(uci -q get gateway.$vlan.ports)"
	local newports

	local name
	local tag

	for port in $oldports; do
		# generate new tagged/untagged syntax
		if [ -z "${port##*t}" ]; then
			# port is tagged
			tag=":t"
		else
			# port is untagged
			tag=":*"
		fi

		# generate new name
		case "$BOARD" in
			tplink,c2600)
				case "${port%%t}" in
					1) name="lan4" ;;
					2) name="lan3" ;;
					3) name="lan2" ;;
					4) name="lan1" ;;
					5) name="wan" ;;
					*) name="" ;;
				esac
				;;
		esac

		# catch empty name (invalid entry or old cpuport)
		[ -z "$name" ] && continue

		# assemble new port name
		newports="$newports ${name}${tag}"
	done

	# cleanup unnecessary spaces
	newports=$(echo "$newports" | xargs)

	uci set gateway.$vlan.ports="$newports"
}

case "$BOARD" in
	tplink,c2600)
		config_load gateway
		config_foreach translate_ports vlan
		;;
esac

uci set gateway.meta.config_version='4'
uci commit gateway

exit 0
