# src/packages/fff/fff-layer3-config/files/etc/uci-defaults/23-network-layer3
#!/bin/sh
# SPDX-License-Identifier: GPL-3.0-or-later

. /lib/functions/system.sh
. /lib/functions/fff/network

BOARD="$(uci get board.model.name)"
. /etc/network.$BOARD

[ -n "$ROUTERMAC" ] || ROUTERMAC=$(get_mac_label)

prefix="fdff::/128"

# Set $prefix::1 as IP
addr1=$(owipcalc "$prefix" add ::1)
# Set $prefix::MAC as IP
addr2=$(owipcalc "$prefix" add "::$(ipMacSuffix "$ROUTERMAC")")
# Set $prefix::EUI64 as IP
addr3=$(owipcalc "$prefix" add "::$(ipEUISuffix "$ROUTERMAC")")

uci batch <<-__EOF__
	set network.fdff='interface'
	set network.fdff.ifname='lo'
	set network.fdff.proto='static'
	set network.fdff.ip6deprecated='1'

	add_list network.fdff.ip6addr="${addr1}/128"
	add_list network.fdff.ip6addr="${addr2}/128"
	add_list network.fdff.ip6addr="${addr3}/128"

	add_list network.client.ip6addr='${addr1}/64'
__EOF__


uci -q commit network
