# src/packages/fff/fff-macnock/files/etc/uci-defaults/52-fff-macnock
#!/bin/sh

if ! grep -q u32 /etc/iproute2/ematch_map;
then
    echo "3 u32" > /etc/iproute2/ematch_map
fi
