# src/packages/fff/fff-batman-adv/files/etc/uci-defaults/93-fff-batman-adv
#!/bin/sh

uci batch <<EOF
  set network.bat0=interface
  set network.bat0.proto=batadv
  set network.bat0.gw_mode='client'
  set network.bat0.gw_sel_class='1'
  set network.bat0.network_coding='0'
  set network.bat0.aggregated_ogms='1'
  set network.bat0.ap_isolation='0'
  set network.bat0.bonding='0'
  set network.bat0.fragmentation='1'
  set network.bat0.orig_interval='1000'
  set network.bat0.distributed_arp_table='1'
  set network.bat0.hop_penalty='30'
EOF

uci commit network

exit 0
