# src/packages/fff/fff-timeserver/files/etc/uci-defaults/51-fff-timeserver
# No initialization here, this is done in configurehood

uci set "system.@system[0].timezone=CET-1CEST,M3.5.0,M10.5.0/3"

uci commit system

exit 0
