# src/packages/fff/fff-timeserver/files/etc/uci-defaults/51-fff-timeserver
#!/bin/sh

. /lib/functions/fff/timeserver
. /etc/community.cfg

setTimeserver "${NTPD_IP}"
