# root_file_system/default/etc/crontabs/root
*/5 * * * *	sh /etc/tincstart.sh
*/5 * * * *	sh /etc/nodewatcher.sh
*/5 * * * *	sh /etc/configurator.sh
0 * * * * sh /etc/configurator.sh sync_hostname

15 01 * * *	rdate -s time.fu-berlin.de > /dev/null

#Enable zapp script if you are running a gateway
#*/1 * * * * /etc/init.d/zapp

