# root_file_system/dir300/etc/crontabs/root
*/5 * * * *	sleep 20 && sh /etc/tincstart.sh
*/5 * * * *	sleep 40 && sh /etc/nodewatcher.sh
*/5 * * * *	sleep 60 && sh /etc/configurator.sh
0 * * * *	sleep 80 && sh /etc/configurator.sh sync_hostname

15 01 * * *	rdate -s time.fu-berlin.de > /dev/null

#Enable zapp script if you are running a gateway
#*/1 * * * * /etc/init.d/zapp

#Reboot dir300 every 3 days at 04:05
#5 4 */3 * *     reboot

