# src/packages/fff/fff-babeld/files/etc/uci-defaults/01-iproute
#!/bin/sh

echo "10 fff" >> /etc/iproute2/rt_tables
