# src/packages/fff/fff-network/files/etc/hotplug.d/net/40-net-rss
#!/bin/sh
# SPDX-License-Identifier: GPL-3.0

case "$(uci get board.model.name)" in
mikrotik,rb5009)
	ethtool -K eth0 receive-hashing on
	;;
solidrun,clearfog-base|\
solidrun,clearfog-pro)
	ethtool -K eth0 receive-hashing on
	ethtool -K eth1 receive-hashing on
	ethtool -K eth2 receive-hashing on
	;;
esac
