# root_file_system/default/etc/hotplug.d/iface/50-fastdstart
#!/bin/sh 
[ "$ACTION" = "ifup" -a "$INTERFACE" = "wan" ] && {
        sleep 3
        sh /etc/fastdstart.sh
}
