# root_file_system/default/etc/hotplug.d/iface/50-tincstart
[ ifup = "$ACTION" -a "$DEVICE" = eth0 ] && {
	sh /etc/tincstart.sh
}
