# root_file_system/default/etc/hotplug.d/iface/50-tincstart
#!/bin/sh 
[ "$ACTION" = "ifup" -a "$INTERFACE" = "wan" ] && {
        sh /etc/tincstart.sh
        sleep 30
        sh /etc/tincstart.sh
}