# root_file_system/default/etc/hotplug.d/iface/51-nodewatcher
[ ifup = "$ACTION" -a "$DEVICE" = eth0 ] && {
        sh /etc/nodewatcher.sh
}
