# root_file_system/default/etc/hotplug.d/iface/51-nodewatcher
#!/bin/sh
[ "$ACTION" = "ifup" -a "$INTERFACE" = "wan" ] && {
        sh /etc/nodewatcher.sh
}