# src/packages/fff/fff-wireless/files/etc/uci-defaults/24b-config-wireless
# SPDX-License-Identifier: GPL-3.0-only
#
# Copyright 2018 Adrian Schmutzler

BOARD="$(uci get board.model.name)"

# antenna gain
case "$BOARD" in
	ubnt,nanostation-loco-m-xw)
		uci set wireless.radio0.antenna_gain=9
		uci commit wireless
		;;
	ubnt,unifi)
		# No support for Unifi AP Pro!
		uci set wireless.radio0.antenna_gain=3
		uci commit wireless
		;;
	ubnt,unifiac-mesh)
		uci set wireless.radio0.antenna_gain=4
		uci set wireless.radio1.antenna_gain=3
		uci commit wireless
		;;
esac

# rssileds: always disable if package is installed
[ -x "/etc/init.d/rssileds" ] && /etc/init.d/rssileds disable

exit 0
