# src/packages/fff/fff-wireless/files/etc/uci-defaults/62-config-wireless
#!/bin/sh
# Copyright 2018 Adrian Schmutzler
# License GPLv3

BOARD="$(uci get board.model.name)"

if [ -s "/etc/wifi.$BOARD" ] ; then
	. "/etc/wifi.$BOARD"
fi
