# src/packages/fff/fff-alfred-monitoring-proxy/files/etc/uci-defaults/52-alfred-master
uci set "alfred.alfred.mode=master"

uci commit alfred

exit 0
