# src/packages/fff/fff-alfred-monitoring-proxy/files/etc/uci-defaults/52-alfred-master
uci -q set "alfred.alfred.mode=master"

uci -q commit alfred

exit 0
