# src/packages/fff/fff-vpn-select/files/etc/hotplug.d/iface/50-vpn-select
#!/bin/sh 
[ "$ACTION" = "ifup" -a "$INTERFACE" = "wan" ] && {
       sleep 3
       /usr/sbin/vpn-select
}

