# src/packages/fff/fff-babel-bird2/files/etc/uci-defaults/30-disable-bird2
/etc/init.d/bird disable
rm -f /etc/init.d/bird

exit 0
