# src/packages/fff/fff-babel-bird2/files/etc/uci-defaults/60-fff-bird-config
# SPDX-License-Identifier: GPL-3.0-only

mv /etc/bird-fff.conf /etc/bird.conf

exit 0
