# src/packages/fff/fff-boardname/files/etc/uci-defaults/15-fff-boardname
BOARD=$(cat /var/sysinfo/board_name)

uci set board.model.name=$BOARD
uci commit board

exit 0
