# src/packages/fff/fff-boardname/files/etc/uci-defaults/15-fff-boardname
BOARD=$(cat /var/sysinfo/board_name)

if uname -a | grep x86_64 > /dev/null; then
	BOARD="x86_64"
fi

uci set board.model.name=$BOARD
uci commit board

exit 0
